# templates/wol-enable.service.j2
[Unit]
Description=Enable Wake-up on LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ wol_ifname }} wol g

[Install]
WantedBy=basic.target

